# systemd/scale-vision.service
[Unit]
Description=scale-vision daemon
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/scale-vision run --config /etc/scale-vision/config.json
Restart=always
RestartSec=2
User=scalevision
Group=scalevision
RuntimeDirectory=scale-vision
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
